# docker/supervisor/cron.d/schedule
* * * * * /usr/local/bin/php /var/www/src/artisan schedule:run >> /var/log/cron.log 2>&1
